# deployments/docker/Xray/supervisord.conf
[supervisord]
nodaemon=true

[program:xray]
command=/usr/bin/xray -c /usr/local/etc/xray/xr.json

[program:au]
command=/usr/local/bin/au -c /usr/local/etc/au/au.json
